# files/fleet-install.service
[Unit]
Requires=fleet.service
After=fleet.service

[Service]
ExecStart=/usr/bin/sleep infinity

[Install]
WantedBy=multi-user.target
